# etc/systemd/system/docker.service.d/override.conf
[Service]
CPUQuota=80%
IOWeight=500
EnvironmentFile=-/etc/default/docker-memory-limits
MemoryMax=${DOCKER_MEMORY_MAX}
MemorySwapMax=${DOCKER_SWAP_MAX}
